# tag-macos/Library/LaunchAgents/net.yewton.sshadd.plist
<?xml version="1.0" encoding="UTF-8"?><!-- -*- mode: xml; -*- -->
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>net.yewton.sshadd</string>
    <key>ProgramArguments</key>
    <array>
      <string>bash</string>
      <string>-c</string>
      <string>~/Library/LaunchAgents/sshadd.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/sshadd.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/sshadd.err</string>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
